# scripts/packaging/hireloom.service
# systemd unit for Linux servers / WSL2.
#
# Install:
#   sudo cp scripts/packaging/hireloom.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hireloom
#
# Tail logs:
#   journalctl -u hireloom -f
#
# Edit before installing:
#   - User=        the local user that owns the repo (NEVER root)
#   - WorkingDirectory=  the absolute path to the repo
#   - EnvironmentFile=   defaults to /etc/hireloom.env; remove if .env in repo

[Unit]
Description=Hireloom dashboard
Documentation=https://github.com/mysticalsin/Hireloom
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/hireloom
EnvironmentFile=-/etc/hireloom.env
EnvironmentFile=/home/YOUR_USERNAME/hireloom/.env
Environment=NODE_ENV=production
Environment=PORT=4747
Environment=HOST=127.0.0.1

ExecStart=/usr/bin/env node apps/web/server.mjs

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/YOUR_USERNAME/hireloom/data /home/YOUR_USERNAME/hireloom/reports /home/YOUR_USERNAME/hireloom/output /home/YOUR_USERNAME/hireloom/config /home/YOUR_USERNAME/hireloom/tmp
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false   # Playwright/V8 needs JIT writes
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=4096
MemoryMax=1500M
CPUQuota=200%

# Standard streams → journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hireloom

[Install]
WantedBy=multi-user.target
