# packaging/career-ops.service
# systemd unit for Linux servers / WSL2.
#
# Install:
#   sudo cp packaging/career-ops.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now career-ops
#
# Tail logs:
#   journalctl -u career-ops -f
#
# Edit before installing:
#   - User=        the local user that owns the repo (NEVER root)
#   - WorkingDirectory=  the absolute path to the repo
#   - EnvironmentFile=   defaults to /etc/career-ops.env; remove if .env in repo

[Unit]
Description=Career-Ops · JobSeeker dashboard
Documentation=https://github.com/santifer/career-ops
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=tony
Group=tony
WorkingDirectory=/home/tony/career-ops
EnvironmentFile=-/etc/career-ops.env
EnvironmentFile=/home/tony/career-ops/.env
Environment=NODE_ENV=production
Environment=PORT=4747
Environment=HOST=127.0.0.1

ExecStart=/usr/bin/env node dashboard-web/server.mjs

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/tony/career-ops/data /home/tony/career-ops/reports /home/tony/career-ops/output /home/tony/career-ops/config /home/tony/career-ops/tmp
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false   # Playwright/V8 needs JIT writes
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=4096
MemoryMax=1500M
CPUQuota=200%

# Standard streams → journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=career-ops

[Install]
WantedBy=multi-user.target
